# wazuh-rules.xml
<!-- Правила реагирования на инциденты ИБ -->

<group name="cron modification-linux">
  <rule id="100007" level="13">
    <if_sid>550, 554</if_sid>
    <field name="file" type="pcre2">^\/var\/spool\/cron\/crontabs</field>
    <description>Cron job has been modified for user "$(uname)".</description>
    <mitre>
      <id>T1053.003</id>
    </mitre>
  </rule>
</group>

<group name="windows, sysmon, sysmon_process-anomalies,">
   <rule id="100000" level="12">
     <if_group>sysmon_event1</if_group>
     <field name="win.eventdata.image">mimikatz.exe</field>
     <description>Sysmon - Suspicious Process - mimikatz.exe on $(win.system.computer) </description>
   </rule>
</group>

<group name="windows_system">

  <rule id="100003" level="13">
    <if_sid>61100</if_sid>
    <field name="win.system.eventID">^7045$</field>
    <description>Создана новая служба Windows $(win.system.eventID) - $(win.eventdata.imagePath) для запуска с корневого пути windows на машине $(win.system.computer)</description>
    <mitre>
      <id>T1543.003</id>
    </mitre>
    <options>no_full_log</options>
    <options>no_email_alert</options>
  </rule>
  
  <rule id="100009" level="6">
    <if_sid>60002</if_sid>
    <field name="win.system.eventID">^5805$</field>
    <description>5805 — ошибка аутентификации NETLOGON сессии на компьютере $(win.system.computer)</description>
    <options>no_full_log</options>
  </rule>

</group>

<group name="lsass_dump">

  <rule id="100015" level="13">
    <if_group>sysmon_event1</if_group>
     <field name="win.eventdata.parentImage" type="pcre2">(?i)powershell.exe</field>
     <field name="win.eventdata.image" type="pcre2">(?i)Windows\\\\(SysWOW64|Temp).+\.exe</field>
    <options>no_full_log</options>
    <description>Было произведено снятие дампа lsass с помощью библиотеки comsvcs.dll на $(win.system.computer)</description>
    <mitre>
      <id>T1003.001</id>
    </mitre>
  </rule>

 <rule id="100016" level="13">
   <if_group>sysmon_event7</if_group>
    <field name="win.eventdata.imageLoaded" type="pcre2">(?i)comsvcs.dll</field>
    <options>no_full_log</options>
    <description> comsvcs.dll на $(win.system.computer) $(win.eventdata.imageLoaded)- $(win.eventdata.ruleName)</description>
  </rule>

</group>


<group name="misp">

  <rule id="3180" level="10">
    <field name="integration">misp</field>
    <match>misp</match>
    <description>MISP Events</description>
    <options>no_full_log</options>
  </rule>
  <rule id="3181" level="5">
    <if_sid>3180</if_sid>
    <field name="misp.error">\.+</field>
    <description>MISP - Ошибка соединения с API</description>
    <options>no_full_log</options>
  </rule>
  <rule id="3182" level="13">
    <field name="misp.category">\.+</field>
    <description>Обнаружение IoC в Threat Intelligence MISP - Category: $(misp.category), Attribute: $(misp.value)</description>
    <mitre>
      <id>T1071.004</id>
    </mitre>
    <options>no_full_log</options>
  </rule>

</group>

<group name="netlogon">

  <rule id="10987" level="3">
    <decoded_as>netlogon</decoded_as>
    <description>Netlogon group</description>
  </rule>

  <rule id="10995" level="13" frequency="40" timeframe="20">
    <if_matched_sid>10987</if_matched_sid>
    <description>Эксплуатация уязвимости ZEROLOGON (CVE-2020-1472)</description>
    <info type="cve">2020-1472</info>
    <mitre>
      <id>T1210</id>
    </mitre>
  </rule>
</group>


<group name="audit,">

  <rule id="80700" level="0">
    <decoded_as>auditd</decoded_as>
    <description>Audit: Messages grouped.</description>
  </rule>

<rule id="100010" level="3">
    <if_sid>80700</if_sid>
    <field name="audit.key">^ping_exec$</field>
    <description>Обнаружено выполнение команды ping</description>
</rule>

  <rule id="80701" level="1">
    <if_sid>80700</if_sid>
    <field name="audit.type">DAEMON_RESUME|DAEMON_START</field>
    <field name="audit.res">success</field>
    <description>Auditd: Daemon Start / Resume.</description>
    <group>gdpr_IV_30.1.g,gpg13_4.14,gpg13_10.1,</group>
  </rule>

  <rule id="20249" level="5">
    <if_sid>80700</if_sid>
    <field name="audit.type">PROCTITLE|PATH|CWD|SYSCALL</field>
    <field name="audit.file.name">\.+</field>
    <field name="audit.command">chmod</field>
    <field name="audit.key">cron_modify</field>
    <description>Auditd: Модификацмя cron-задачи</description>
  </rule>

</group>

<!-- Этапы интеграции систем -->


  <integration>
    <name>custom-w2thive</name>
    <hook_url>url_address_thehive:9000</hook_url>
    <api_key>*****api-key*****</api_key>
    <alert_format>json</alert_format>
  </integration>

  <integration>
    <name>custom-w2thive5</name>
    <hook_url>url_address_thehive:9000</hook_url>
    <api_key>*****api-key*****</api_key>
    <alert_format>json</alert_format>
  </integration>

<integration>
        <name>custom-telegram</name>
        <level>4</level>
        <hook_url>https://api.telegram.org/bot*****api-key*****:*****api-key*****/sendMessage</hook_url>
        <alert_format>json</alert_format>
    </integration>

  <integration>
    <name>custom-misp</name>
    <group>sysmon_event1,sysmon_event3,sysmon_event6,sysmon_event7,sysmon_event_15,sysmon_event_22</group>
    <alert_format>json</alert_format>
  </integration>


  <integration>
    <name>custom-misp2</name>
    <alert_format>json</alert_format>
  </integration>

<!-- Источники данных linux -->

ossec_config>
  <localfile>
    <log_format>syslog</log_format>
    <location>/var/ossec/logs/active-responses.log</location>
  </localfile>

  <localfile>
    <log_format>syslog</log_format>
    <location>/var/log/auth.log</location>
  </localfile>

  <localfile>
    <log_format>syslog</log_format>
    <location>/var/log/syslog</location>
  </localfile>

  <localfile>
    <log_format>syslog</log_format>
    <location>/var/log/dpkg.log</location>
  </localfile>

  <localfile>
    <log_format>syslog</log_format>
    <location>/var/log/kern.log</location>
  </localfile>

  <localfile>
    <log_format>audit</log_format>
    <location>/var/log/audit/audit.log</location>
  </localfile>

  <localfile>
    <log_format>syslog</log_format>
    <location>/var/log/suricata/eve.json</location>
  </localfile>

  <localfile>
    <log_format>syslog</log_format>
    <location>/var/log/suricata/fast.json</location>
  </localfile>

</ossec_config>

<!-- Источники данных Windows -->

<!-- Log analysis -->
  <localfile>
    <location>Application</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <localfile>
    <location>Security</location>
    <log_format>eventchannel</log_format>
    <query>Event/System[EventID != 5145 and EventID != 5156 and EventID != 5447 and
      EventID != 4656 and EventID != 4658 and EventID != 4663 and EventID != 4660 and
      EventID != 4670 and EventID != 4690 and EventID != 4703 and EventID != 4907 and
      EventID != 5152 and EventID != 5157]</query>
  </localfile>

  <localfile>
    <location>System</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <localfile> 
    <location>Microsoft-Windows-Sysmon/Operational</location> 
    <log_format>eventchannel</log_format> 
  </localfile>

  <localfile>
    <location>active-response\active-responses.log</location>
    <log_format>syslog</log_format>
  </localfile>

